# deploy/systemd/platform-api.service
[Unit]
Description=SGX OneBox Platform API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sgx-onebox
EnvironmentFile=/opt/sgx-onebox/config/platform.env
Environment=PLATFORM_RUNTIME_BUNDLES_DIR=/opt/sgx-onebox/runtime-bundles
ExecStart=/opt/sgx-onebox/bin/platform-api -static /opt/sgx-onebox/web -json /opt/sgx-onebox/data/platform-data.json
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
